# premdev/infra/premdev-bot.service
[Unit]
Description=PremDev Telegram Admin Bot
After=network-online.target docker.service
Wants=network-online.target
# If the bot somehow crashes during a docker restart of the stack, give it
# room to come back rather than entering a tight restart loop.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/env node /usr/local/sbin/premdev-bot.mjs
Restart=on-failure
RestartSec=10
# The bot reads /opt/premdev/.env directly (NOT through EnvironmentFile=, so
# ADMIN_TELEGRAM_ID changes via /setenv take effect after `systemctl restart`).
WorkingDirectory=/opt/premdev
StandardOutput=journal
StandardError=journal
SyslogIdentifier=premdev-bot

# Light hardening — bot needs Docker socket + .env write so we can't lock
# things down too far, but block obvious unrelated attack surfaces.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target
